# bin/gpio.service
[Unit]
 Description=Read GPIO Service
 After=multi-user.target

 [Service]
 User=pi
 Group=pi
 Type=idle
 #ExecStart=/usr/bin/gpio-watch -s /home/pi/bin/piSchoolBell/gpio-scripts -e switch 7 8 9 10
 ExecStart=/usr/bin/gpio-watch -s /home/pi/bin/piSchoolBell/gpio-scripts -e switch 7 8 9 10 -l /home/pi/bin/piSchoolBell/gpio-watch.log > /dev/null 2>&1

 [Install]
 WantedBy=multi-user.target
 